# cookbook/templates/heimdall-lite.service.erb
#
# /etc/systemd/system/heimdall-list.service
#
# Generated by Chef for <%= node['fqdn'] %>
#


[Unit]
Description=MITRE Heimdall Lite Server
After=network.target

[Service]
Type=simple
ExecStart=/usr/bin/heimdall-lite <%= @port -%>
Restart=always
StandardOutput=syslog
StandardError=syslog
SyslogIdentifier=heimdall-lite
User=root
Environment=NODE_ENV=production

[Install]
WantedBy=multi-user.target
